# base/overlay/usr/lib/systemd/system/btrfs-scrub@.timer
[Unit]
Description=Monthly Btrfs scrub timer

[Timer]
OnCalendar=monthly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
